# systemd/troop-backup.timer
[Unit]
Description=Monkey Troop Database Backup Timer
Documentation=https://github.com/100monkeys-ai/monkey-troop
Requires=troop-backup.service

[Timer]
# Run daily at 2:00 AM
OnCalendar=daily
OnCalendar=02:00

# If system was off at scheduled time, run on next boot
Persistent=true

# Randomize execution time by up to 15 minutes to avoid load spikes
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
